# contrib/systemd/kportwatch.service
# KPortWatch — Systemd user service
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp contrib/systemd/kportwatch.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now kportwatch
#
# Status:
#   systemctl --user status kportwatch
#
# Logs:
#   journalctl --user -u kportwatch -f
#
# The service runs under the user's XDG_RUNTIME_DIR so the data file
# is only readable by the owner (mode 0700 directory).

[Unit]
Description=KPortWatch Network Security Monitor
Documentation=https://github.com/user/kportwatch
After=network.target

[Service]
Type=simple
ExecStart=%h/KPortWatch/.venv/bin/python -m backend.kportwatch_daemon --foreground
WorkingDirectory=%h/KPortWatch
Restart=on-failure
RestartSec=5

# Logging → journal (view with journalctl --user -u kportwatch)
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/kportwatch %h/.local/share/kportwatch %t
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
