# systemd/netsentry.service
# NetSentry — Systemd user service
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp contrib/systemd/netsentry.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now netsentry
#
# Status:
#   systemctl --user status netsentry
#
# Logs:
#   journalctl --user -u netsentry -f
#
# The service runs under the user's XDG_RUNTIME_DIR so the data file
# is only readable by the owner (mode 0700 directory).

[Unit]
Description=NetSentry Network Security Monitor
Documentation=https://github.com/user/netsentry
After=network.target

[Service]
Type=simple
ExecStart=%h/NetSentry/.venv/bin/python -m backend.netsentry_daemon --foreground
WorkingDirectory=%h/NetSentry
Restart=on-failure
RestartSec=5

# Logging → journal (view with journalctl --user -u netsentry)
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/netsentry %t
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
